# ttdash.service
# sudo ln -s $(pwd)/ttdash.service /lib/systemd/system/ttdash.service
[Unit]
Description=TrainTrack Dashboard

[Service]
Type=simple
ExecStart=/home/pi/ttdash --auto-update --debug-port=4242 --purpleair-credentials-file /home/pi/src/ttdash/purpleair-credentials.json --mta-bustime-credentials-file /home/pi/src/ttdash/bustime.creds

[Install]
WantedBy=multi-user.target
